# dot_local/etc/profile.d/50-fzf
FZF_DEFAULT_OPTS=""
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --multi"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --border rounded"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --layout reverse"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --info hidden"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --height 13"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --scroll-off 3"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --prompt '⟩ '"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --pointer '▶'"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --marker '●'"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color fg:-1:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color bg:-1"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color preview-fg:-1:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color preview-bg:-1"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color hl:red:bold"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color fg+:yellow:bold"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color bg+:-1"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color gutter:-1"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color hl+:red:bold"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color query:-1:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color disabled:black:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color info:black:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color border:black:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color prompt:cyan:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color pointer:yellow:bold"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color marker:blue:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color spinner:magenta:regular"
FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS --color header:-1:bold"
export FZF_DEFAULT_OPTS="$FZF_DEFAULT_OPTS"
