# dot_local/etc/profile.d/50-python
export PYTHONDONTWRITEBYTECODE=1
export PYTHONPATH="$PYTHONPATH:$HOME/.local/lib/python3/custom-packages"
export PYTHONSTARTUP="$XDG_CONFIG_HOME/python/pythonrc"
